# systemd/rpi_temp.timer
[Unit]
Description=RPi temperature log every 2 minutes

[Timer]
OnCalendar=*-*-* *:00/2:00
Persistent=true

[Install]
WantedBy=timers.target
